# sakura/roles/systemd_service/templates/unit.service.j2
# {{ ansible_managed }}

[Unit]
{% for key, value in unit.items() %}
{{ key }} = {{ value }}
{% endfor %}

[Service]
{% for key, value in service.items() %}
{{ key }} = {{ value }}
{% endfor %}

[Install]
{% for key, value in install.items() %}
{{ key }} = {{ value }}
{% endfor %}
