# sakura/roles/systemd_service/templates/unit.service.j2
# {{ ansible_managed }}

[Unit]
Description = {{ description }}

[Service]
ExecStart = {{ execstart }}
Restart = always
Type = {{ type }}

[Install]
WantedBy = multi-user.target
