# sakura/docker-compose/nextcloud0/crontab
40 3,18,33,48 * * * * php -f /var/www/html/cron.php
