# _conoha/roles/wiki/templates/supervisord.conf.j2
[program:gollum]
command=make start
autorestart=true
user={{ wiki_user }}
directory=/home/{{ wiki_user }}/gollum-site
