# conoha/docker-compose/mastodon/crontab
45 3,15 * * * rake mastodon:daily
15 8 * * 1 rake mastodon:media:remove_remote
